# server/deploy/plantcart-backup.timer
[Unit]
Description=PlantCart nightly DB backup (03:30, keep 14)

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
