# server/deploy/thincart-backup.service
# ThinCart nightly DB backup (systemd user unit — triggered by the .timer)
# The purchase history IS the intelligence; losing it resets the app.
[Unit]
Description=ThinCart SQLite nightly backup

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'mkdir -p %h/backups/thincart && sqlite3 %h/Project/thincart/server/data/thincart.db ".backup %h/backups/thincart/thincart-$(date +%%F).db" && ls -1t %h/backups/thincart/thincart-*.db | tail -n +15 | xargs -r rm'
